# docker/root/etc/services.d/apache2/run
#!/usr/bin/with-contenv bash

echo "Starting Web Server"

exec \
	/usr/sbin/apachectl -D FOREGROUND